# runtime/agent_orchestrator/launchd/com.kelvin.lottery.copilot-daemon.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key><string>com.kelvin.lottery.copilot-daemon</string>
    <key>ProgramArguments</key>
    <array>
        <string>/usr/bin/python3</string>
        <string>/Users/kelvin/Kelvin-WorkSpace/LotteryNew/orchestrator/copilot_daemon.py</string>
        <string>--poll-seconds</string>
        <string>10</string>
    </array>
    <key>WorkingDirectory</key><string>/Users/kelvin/Kelvin-WorkSpace/LotteryNew</string>
    <key>RunAtLoad</key><true/>
    <key>KeepAlive</key><true/>
    <key>EnvironmentVariables</key>
    <dict>
        <key>HOME</key><string>/Users/kelvin</string>
        <key>PATH</key><string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
        <key>GH_CONFIG_DIR</key><string>/Users/kelvin/.config/gh</string>
    </dict>
    <key>StandardOutPath</key><string>/Users/kelvin/Kelvin-WorkSpace/LotteryNew/runtime/agent_orchestrator/logs/copilot-daemon.launchd.log</string>
    <key>StandardErrorPath</key><string>/Users/kelvin/Kelvin-WorkSpace/LotteryNew/runtime/agent_orchestrator/logs/copilot-daemon.launchd.err</string>
</dict>
</plist>
